# deploy/gunicorn.service
[Unit]
Description=Network Monitor — Gunicorn
After=network.target postgresql.service

[Service]
User=netmon
Group=netmon
WorkingDirectory=/opt/network_dashboard
EnvironmentFile=/opt/network_dashboard/.env
ExecStart=/opt/network_dashboard/venv/bin/gunicorn \
    --workers 3 \
    --bind 0.0.0.0:8000 \
    network_dashboard.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
